# _orbicon-dev/orbicon/_non-bank/stock-exch-cron-old/zse.crontab
*/15 * * * * php -f /PATH_TO_DOCUMENT_ROOT/orbicon/modules/stock-exch-cron/cron.zse.php