# _orbicon-dev/orbicon/modules/synchronization/sync.crontab
1 18 * * * php -f /PATH_TO_DOCUMENT_ROOT/orbicon/modules/synchronization/cron.sync.php